# devops/crontabs/notify.cron
0 17 * * * openshiksha_app /home/openshiksha_app/openshiksha/pylon/scripts/notify-due.sh >> /home/openshiksha_app/openshiksha/devops/cron.log 2>&1

